# 1.ini
[program:wxbot] ; 进程名字

command=  /root/wxbot/cmd/cmd;启令 

user=root ; 启动进程的用户 

autostart=true                ; 随着supervisord的启动而启动

autorestart=true              ; 自动重启

startsecs=4                   ; 启动10秒后没有异常退出，就表示进程正常启动了，默认为1

stderr_logfile=/var/log/wxbot.err.log  ;错误日志

stdout_logfile=/var/log/wxbot.out.log ;info 日志